# Crypto/300_AlphaComplex1/build/deploy/cookbooks/alpha-complex-1/files/etc/xinetd.d/alpha-complex-1
service alpha-complex-1
{
    socket_type     = stream
    type            = UNLISTED
    user            = alpha-complex-1
    port            = 31337
    server          = /home/alpha-complex-1/server.sh
    wait            = no
}

